# roles/apply_health_metrics_hooks/files/systemd/health-hook@.timer
[Unit]
Description=Health-dashboard textfile-collector timer: %i

[Timer]
OnBootSec=30s
# OnUnitActiveSec is set per-instance via a drop-in (see role tasks).
Unit=health-hook@%i.service

[Install]
WantedBy=timers.target
