# roles/apply_health_metrics_hooks/files/systemd/health-hook@.service
[Unit]
Description=Health-dashboard textfile-collector hook: %i
After=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/flax/sbin/health_hooks/%i.sh
User=root
Group=root
