# roles/apply_kea/templates/kea-dhcp4.service.j2
[Unit]
Description=ISC Kea DHCPv4 (all braintree v4 vlans)
After=docker.service network-online.target flax-stack.service
Wants=network-online.target
Requires=docker.service flax-stack.service

[Service]
Type=simple
EnvironmentFile=/etc/flax/kea/kea.env

# Render the password into the JSON config from the env file. Kea 2.6
# does NOT natively interpolate env vars in JSON; envsubst handles it
# at unit start. Mode 0640 + chown to the kea container user (uid 100,
# gid 101 in the cloudsmith image) so the in-container daemon can read
# it (the config still excludes "other" so non-kea processes on the
# host can't read the embedded Postgres password).
ExecStartPre=/bin/bash -c 'umask 027 && envsubst < /etc/flax/kea/kea-dhcp4.conf.template > /etc/flax/kea/kea-dhcp4.conf && chown 100:101 /etc/flax/kea/kea-dhcp4.conf && chown 100:101 /var/run/kea'

# Override --entrypoint to run kea-dhcp4 directly. The cloudsmith image's
# CMD is supervisord which also starts kea-ctrl-agent (needs its own
# /etc/kea/kea-ctrl-agent.conf that we don't provide). We don't need the
# REST control agent for Plan 5; skipping it via direct entrypoint
# sidesteps the missing-config error.
ExecStart=/usr/bin/docker run --rm --name kea-dhcp4 \
    --network host \
    --cap-add NET_ADMIN --cap-add NET_RAW --cap-add NET_BIND_SERVICE \
    -v /etc/flax/kea/kea-dhcp4.conf:/etc/kea/kea-dhcp4.conf:ro \
    -v /var/run/kea:/var/run/kea \
    --env-file /etc/flax/kea/kea.env \
    --entrypoint /usr/sbin/kea-dhcp4 \
    docker.cloudsmith.io/isc/docker/kea-dhcp4:${KEA_IMAGE_TAG} \
    -c /etc/kea/kea-dhcp4.conf

ExecStop=/usr/bin/docker stop kea-dhcp4

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
